# scripts/com.jobagent.run.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobagent.run</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/muhammadeltahir/Projects/Job Search/scripts/run_scheduled.sh</string>
    </array>

    <!-- Wall-clock trigger: attempt daily at 9:00 AM. Far more reliable than
         StartInterval — if the Mac is asleep at 9am, launchd runs it at the next wake.
         `run --if-due` then skips unless >=47h since the last success, so a real run
         happens about every other day without ever silently missing. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/muhammadeltahir/Projects/Job Search</string>

    <key>StandardOutPath</key>
    <string>/Users/muhammadeltahir/Projects/Job Search/data/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/muhammadeltahir/Projects/Job Search/data/launchd.err.log</string>
</dict>
</plist>
